# deploy/vps/systemd/openalice-alpaca-smoke-week.service
[Unit]
Description=Run OpenAlice Alpaca paper smoke test
After=docker.service

[Service]
Type=oneshot
WorkingDirectory=/srv/bots/openalice-paper-eval
ExecStart=/bin/sh /srv/bots/openalice-paper-eval/scripts/run_alpaca_smoke_week.sh
